# image_debian_stretch/etc/rc.local
#!/bin/bash
if [ -f "/root/.boot" ] ;  then
# /root/.boot &
fi

