# linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/libopencl1.list
/usr/lib/libOpenCL.so.1
/usr/share/ti/opencl/profiling/classes/__init__.py
/etc/ti-mctd/ti_mctd_config.json
/usr/share/ti/opencl/cpu.h
/usr/share/ti/opencl/dsp.h
/usr/share/ti/opencl/dsp_edmamgr.h
/usr/share/ti/opencl/profiling/classes/
/lib/systemd/system/ti-mct-daemon.service
/usr/lib/libOpenCL.so.1.1.14
/usr/share/ti/opencl/profiling/classes/event.py
/usr/share/ti/opencl/profiling/
/usr/share/ti/opencl/clc.h
/usr/share/ti/opencl/profiling/classes/kernel.py
/usr/share/ti/opencl/profiling/classes/graph.py
/usr/share/ti/opencl/dsp_c.h
/usr/share/ti/opencl/profiling/classes/event_list.py
/usr/share/ti/opencl/profiling/classes/table.py
/usr/share/ti/opencl/profiling/oclaet.py
/usr/share/ti/opencl/profiling/oclaet.sh
/etc/ti-mctd/
/usr/share/ti/opencl/profiling/classes/core.py
/usr/share/ti/opencl/
