# linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/run-postinsts.list
/lib/systemd/system/run-postinsts.service
/etc/init.d/run-postinsts
/etc/init.d/
/usr/sbin/run-postinsts
