# linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/ti-ipc.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/tiipclad-daemon.sh" ]; then
	/etc/init.d/tiipclad-daemon.sh stop || :
fi
# End section update-rc.d
