# linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/ti-sgx-ddk-um.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/rc.pvr" ]; then
	/etc/init.d/rc.pvr stop || :
fi
# End section update-rc.d
