# linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/ti-sgx-ddk-um.list
/usr/lib/libpvr_wlegl.so.1.14.3699939
/usr/lib/libGLES_CM.so.1.14.3699939
/usr/lib/libusc.so.1
/usr/lib/libsrv_init.so.1
/usr/lib/libpvr2d.so.1
/usr/lib/libpvrDRMWSEGL.so.1
/usr/lib/libIMGegl.so.1.14.3699939
/etc/powervr.ini
/usr/lib/libpvr_wlegl.so.1
/usr/lib/libusc.so.1.14.3699939
/usr/lib/libEGL.so.1.14.3699939
/usr/lib/libsrv_init.so.1.14.3699939
/usr/bin/glsltest1_fragshaderA.txt
/usr/bin/eglinfo
/usr/lib/libglslcompiler.so.1.14.3699939
/usr/lib/libsrv_um.so.1
/usr/lib/libglslcompiler.so.1
/usr/lib/libpvr2d.so.1.14.3699939
/usr/lib/libpvrws_WAYLAND.so.1
/usr/bin/glsltest1_fragshaderB.txt
/usr/lib/libPVRScopeServices.so.1
/usr/lib/libIMGegl.so.1
/usr/lib/libGLES_CM.so.1
/usr/lib/libsrv_um.so.1.14.3699939
/usr/lib/libPVRScopeServices.so.1.14.3699939
/usr/lib/libEGL.so.1
/usr/lib/gbm/
/usr/lib/gbm/gbm_pvr.so
/usr/lib/libdbm.so.1
/usr/lib/libpvrDRMWSEGL.so.1.14.3699939
/usr/lib/libdbm.so.1.14.3699939
/usr/bin/gles1test1
/usr/lib/libpvrws_WAYLAND.so.1.14.3699939
/usr/bin/glsltest1_vertshader.txt
/usr/lib/libpvrGBMWSEGL.so.1
/usr/lib/libpvrGBMWSEGL.so.1.14.3699939
/usr/bin/pvrsrvctl
/etc/init.d/rc.pvr
/usr/bin/gles2test1
/usr/lib/libGLESv2.so.1
/usr/lib/libGLESv2.so.2
/usr/lib/libGLESv2.so.1.14.3699939
